# rpmbuild/SOURCES/omp2_base/etc/profile.d/init_java_home.sh
#!/bin/bash

export JAVA_HOME=/usr/share/java/jdk1.7.0_45
export PATH=$JAVA_HOME/bin:$PATH
