# rpmbuild/SOURCES/omp2_base/etc/profile.d/storm-init.sh
#!/bin/bash
export STORM_HOME=/usr/local/storm
export PATH=$PATH:$STORM_HOME/bin