# rpmbuild/SOURCES/omp2_base/service/e_drpc/run
#!/bin/sh
exec /usr/local/storm/bin/storm_run drpc 
