# rpmbuild/SOURCES/omp2_base/service/d_neo4j/run
#!/bin/sh
exec 2>&1
exec /usr/share/neo4j-community-1.9.5/bin/neo4j start-no-wait
