# rpmbuild/SOURCES/omp2_base/service/f_storm-ui/run
#!/bin/sh
exec 2>&1
exec /usr/local/storm/bin/storm_run ui 
